# files/etc/crontabs/root
*/1 * * * * /usr/local/bin/check_connectivity.sh