# files/etc/uci-defaults/89-tollgate-ecash-test
#!/bin/sh

# Create test file to verify permissions
mkdir -p /etc/tollgate/ecash
touch /etc/tollgate/ecash/test-file.txt
echo "Test content $(date)" > /etc/tollgate/ecash/test-file.txt
chmod 777 -R /etc/tollgate/ecash

# Make a link in /tmp
ln -sf /etc/tollgate/ecash /tmp/tollgate-ecash

# Add debug info
echo "File test completed at $(date)" > /tmp/ecash-test-result.log
ls -la /etc/tollgate/ecash >> /tmp/ecash-test-result.log

exit 0 