# files/etc/hotplug.d/iface/95-tollgate-restart
#!/bin/sh

# This script handles the one-time restart of the tollgate-wrt service
# after the WAN interface comes up and internet connectivity is confirmed.

# --- DEBUG: Log all environment variables passed to the script ---
logger "DEBUG tollgate-wrt hotplug: ACTION='$ACTION', INTERFACE='$INTERFACE', DEVICE='$DEVICE'"
# --- END DEBUG ---

[ "$ACTION" = "ifup" ] || exit 0

# Check if reseller_mode is enabled. If so, exit immediately.
# The wireless_gateway_manager is responsible for all wwan interface handling.
reseller_mode=$(jsonfilter -e '@.reseller_mode' < /etc/tollgate/config.json 2>/dev/null)
if [ "$reseller_mode" = "true" ]; then
    logger "DEBUG tollgate-wrt hotplug: Reseller mode enabled, 95-tollgate-restart exiting."
    exit 0
fi

# Get the actual device name for the 'wan' logical interface
WAN_DEVICE=$(uci get network.wan.device 2>/dev/null)

# Check if the interface that came up is the one associated with 'wan' or 'wwan'
# This handles cases where wan might use a wireless interface like 'wwan'
if [ -n "$WAN_DEVICE" ] && [ "$DEVICE" = "$WAN_DEVICE" ]; then
    # This is the wan interface coming up
    logger "DEBUG tollgate-wrt hotplug: Matched device '$DEVICE' to UCI wan.device '$WAN_DEVICE'"
    :
elif [ "$INTERFACE" = "wan" ] || [ "$INTERFACE" = "wwan" ]; then
    # Fallback to checking INTERFACE if DEVICE didn't match
    # This covers standard configurations and wireless configurations
    logger "DEBUG tollgate-wrt hotplug: Matched INTERFACE '$INTERFACE'"
    :
else
    # Not the wan interface we're interested in
    logger "DEBUG tollgate-wrt hotplug: No match. ACTION='$ACTION', INTERFACE='$INTERFACE', DEVICE='$DEVICE', UCI wan.device='$WAN_DEVICE'"
    exit 0
fi

FLAG_FILE="/tmp/tollgate_initial_restart_done"

# Check if the restart has already been performed
if [ -f "$FLAG_FILE" ]; then
    # Already restarted once, nothing to do
    logger "DEBUG tollgate-wrt hotplug: Flag file '$FLAG_FILE' found, exiting."
    exit 0
fi

# Function to check internet connectivity
check_connectivity() {
    ping -c 1 -W 5 8.8.8.8 >/dev/null 2>&1
    return $?
}

# Wait a bit for the interface to be fully configured
logger "DEBUG tollgate-wrt hotplug: Waiting 5 seconds for interface to settle..."
sleep 5

# Check for connectivity
logger "DEBUG tollgate-wrt hotplug: Checking for internet connectivity..."
if check_connectivity; then
    logger "DEBUG tollgate-wrt hotplug: Internet connectivity confirmed."
    # Restart the tollgate-wrt service
    logger "tollgate-wrt: Restarting service after initial WAN up and connectivity check."
    /etc/init.d/tollgate-wrt restart
    
    # Create the flag file to indicate the restart has been done
    touch "$FLAG_FILE"
    
    # Log the action
    logger "tollgate-wrt: Service restarted and flag file created."
else
    # No connectivity yet, exit silently
    # The script will run again on the next wan ifup event if needed
    logger "tollgate-wrt: WAN up but no connectivity to 8.8.8.8, skipping restart for now."
fi