# files/etc/hotplug.d/iface/96-tollgate-scan
#!/bin/sh

# This script triggers a network scan in the tollgate-wrt service
# when a WAN interface comes up.

# Log script execution start
logger -t "tollgate-hotplug-scan" "Executing 96-tollgate-scan script. ACTION: $ACTION, INTERFACE: $INTERFACE"

[ "$ACTION" = "ifup" ] || exit 0

# This script now runs in all modes. The tollgate-cli will handle
# the distinction between a reseller scan and a standard scan.

# Check if the interface is 'wan' or 'wwan'
if [ "$INTERFACE" = "wan" ] || [ "$INTERFACE" = "wwan" ]; then
    logger -t "tollgate-hotplug-scan" "Matched INTERFACE '$INTERFACE'. Finding real device name..."

    # The 'ubus' command is the standard way in OpenWrt to get the real device name
    # associated with a logical interface like 'wwan'. We poll briefly in case the
    # interface is not immediately available in ubus after the ifup event.
    real_device=""
    attempts=0
    while [ $attempts -lt 10 ]; do
        real_device=$(ubus -S call network.interface.$INTERFACE status | jsonfilter -e '@.l3_device')
        [ -n "$real_device" ] && break
        sleep 1
        attempts=$((attempts + 1))
    done

    if [ -z "$real_device" ]; then
        logger -t "tollgate-hotplug-scan" "Could not determine real device for '$INTERFACE'. Aborting."
        exit 1
    fi

    logger -t "tollgate-hotplug-scan" "Real device for '$INTERFACE' is '$real_device'. Waiting for default route..."

    # Wait for the default route to be established on the real device
    attempts=0
    max_attempts=30 # Wait for up to 30 seconds
    route_found=0
    while [ $attempts -lt $max_attempts ]; do
        if ip route show default | grep -q "dev $real_device"; then
            logger -t "tollgate-hotplug-scan" "Default route for '$real_device' found after $attempts seconds."
            route_found=1
            break
        fi
        sleep 1
        attempts=$((attempts + 1))
    done

    if [ $route_found -eq 0 ]; then
        logger -t "tollgate-hotplug-scan" "Timed out waiting for default route on '$real_device'. Aborting scan."
        exit 1
    fi
    
    # Trigger a network scan via the CLI
    logger -t "tollgate-hotplug-scan" "Executing tollgate network scan..."
    /usr/bin/tollgate network scan
    SCAN_EXIT_CODE=$?
    logger -t "tollgate-hotplug-scan" "tollgate network scan finished with exit code: $SCAN_EXIT_CODE"
else
    logger -t "tollgate-hotplug-scan" "No match for INTERFACE '$INTERFACE'."
    exit 0
fi

logger -t "tollgate-hotplug-scan" "Finished 96-tollgate-scan script."