# files/etc/hotplug.d/iface/96-tollgate-scan
#!/bin/sh

# This script triggers a network scan in the tollgate-wrt service
# when a WAN interface comes up.

# Log script execution start
logger -t "tollgate-hotplug-scan" "Executing 96-tollgate-scan script. ACTION: $ACTION, INTERFACE: $INTERFACE"

[ "$ACTION" = "ifup" ] || exit 0

# Check if the interface is 'wan' or 'wwan'
if [ "$INTERFACE" = "wan" ] || [ "$INTERFACE" = "wwan" ]; then
    logger -t "tollgate-hotplug-scan" "Matched INTERFACE '$INTERFACE'. Waiting for default route..."

    # Wait for the default route to be established on the interface
    attempts=0
    max_attempts=30 # Wait for up to 30 seconds
    route_found=0
    while [ $attempts -lt $max_attempts ]; do
        # Check for a default route via the specified interface
        if ip route show default | grep -q "dev $INTERFACE"; then
            logger -t "tollgate-hotplug-scan" "Default route for '$INTERFACE' found after $attempts seconds."
            route_found=1
            break
        fi
        sleep 1
        attempts=$((attempts + 1))
    done

    if [ $route_found -eq 0 ]; then
        logger -t "tollgate-hotplug-scan" "Timed out waiting for default route on '$INTERFACE'. Aborting scan."
        exit 1
    fi
    
    # Trigger a network scan via the CLI
    logger -t "tollgate-hotplug-scan" "Executing tollgate-cli network scan."
    /usr/bin/tollgate-cli network scan
    logger -t "tollgate-hotplug-scan" "tollgate-cli network scan finished."
else
    logger -t "tollgate-hotplug-scan" "No match for INTERFACE '$INTERFACE'."
    exit 0
fi

logger -t "tollgate-hotplug-scan" "Finished 96-tollgate-scan script."